# scripts/cloud-image/systemd/xelora-firstboot.service
[Unit]
Description=Xelora first-boot initialization (runs once)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target
Before=xelora.service
ConditionPathExists=/usr/local/sbin/xelora-firstboot.sh
ConditionPathExists=!/opt/Xelora/.firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/xelora-firstboot.sh
RemainAfterExit=no
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
